# dist/etc/supervisor/conf.d/10-apache2.conf
[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND
autostart=true
autorestart=true
startretries=1
startsecs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2.error.log
user=root
killasgroup=true
stopasgroup=true
