# dist/etc/supervisor/conf.d/10-apache2.conf
[program:nginx]
command=/usr/sbin/apache2 -DFOREGROUND
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2.error.log
