# dist/etc/supervisor/conf.d/10-puma.conf
[program:puma]
command=bundle exec puma
directory=/var/www/app
user=www-data
stdout_logfile=/var/log/supervisor/puma.log
stderr_logfile=/var/log/supervisor/puma.error.log
