# dist/etc/supervisor/conf.d/10-shibd.conf
[program:shibd]
command=/usr/sbin/shibd -F -w30
user=_shibd
stdout_logfile=/var/log/supervisor/shibd.log
stderr_logfile=/var/log/supervisor/shibd.error.log
