# dist/etc/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.error.log

[program:puma]
command=/insert/correct/path/to/puma
user=www-data
stdout_logfile=/var/log/supervisor/puma.log
stderr_logfile=/var/log/supervisor/puma.error.log

[program:shibd]
command=/usr/sbin/shibd -F
user=_shibd
stdout_logfile=/var/log/supervisor/shibd.log
stderr_logfile=/var/log/supervisor/shibd.error.log

[fcgi-program:shibauthorizer]
command=/usr/lib/x86_64-linux-gnu/shibboleth/shibauthorizer
socket=unix:///run/shibboleth/shibauthorizer.sock
socket_owner=shibd:www-data
socket_mode=0660
user=_shibd
stdout_logfile=/var/log/supervisor/shibauthorizer.log
stderr_logfile=/var/log/supervisor/shibauthorizer.error.log

[fcgi-program:shibresponder]
command=/usr/lib/x86_64-linux-gnu/shibboleth/shibresponder
socket=unix:///run/shibboleth/shibresponder.sock
socket_owner=shibd:www-data
socket_mode=0660
user=_shibd
stdout_logfile=/var/log/supervisor/shibresponder.log
stderr_logfile=/var/log/supervisor/shibresponder.error.log
